# systemd/bob.service
[Unit]
Description=BoB Assistant (jour 1 – TTS demo)
After=network-online.target sound.target

[Service]
Type=simple
User=bob
WorkingDirectory=/home/bob/projects/bob
Environment=BOB_CONFIG=/home/bob/projects/bob/config/bob.yaml
ExecStart=/home/bob/projects/bob/target/release/bob say "BoB est opérationnel."
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
